# focus/systemd/focus.service
[Unit]
Description=FOCUS Productivity Agent Daemon
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/focus
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/focus/.venv/bin/python %h/focus/focus_daemon.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
